# rtltcp_437.conf
#!upstart

description "Upstart script for 437 MHz rtl_tcp"
 
start on starting websdr
stop on shutdown or stopped websdr
respawn

script
    echo "Starting 437 MHz rtl_tcp.."
    exec rtl_tcp -d 4 -f 437000000 -s 2048000 -g 30 -a 127.0.0.1 -p 5905
end script
