# rtltcp_432.conf
#!upstart

description "Upstart script for 432 MHz rtl_tcp"
 
start on starting websdr
stop on shutdown or stopped websdr
respawn

script
    echo "Starting 432 MHz rtl_tcp.."
    exec rtl_tcp -d 3 -f 433700000 -s 2048000 -g 40 -a 127.0.0.1 -p 5902
end script
