# rtltcp_144.conf
#!upstart

description "Upstart script for 144 MHz rtl_tcp"
 
start on starting websdr
stop on shutdown or stopped websdr
respawn

script
    echo "Starting 144 MHz rtl_tcp.."
    exec rtl_tcp -d 2 -f 145000000 -s 2048000 -g 35 -a 127.0.0.1 -p 5901
end script
