# rtltcp_146.conf
#!upstart

description "Upstart script for 146 MHz rtl_tcp"
 
start on starting websdr
stop on shutdown or stopped websdr
respawn

script
    echo "Starting 146 MHz rtl_tcp.."
    exec rtl_tcp -d 5 -f 146000000 -s 2048000 -g 30 -a 127.0.0.1 -p 5906 -b 3 -n 1000
end script
