# rtltcp_432.conf
#!upstart

description "Upstart script for 432 MHz rtl_tcp"
 
start on starting websdr
stop on shutdown or stopped websdr
respawn
nice -1

script
    echo "Starting 432 MHz rtl_tcp.."
    exec rtl_tcp -d 4 -f 433700000 -s 2048000 -g 30 -a 127.0.0.1 -p 5903 -b 3 -n 1000
end script
